# docker/media_server/supervisord.conf
[supervisord]
nodaemon=true

[program:plex]
command=/home/plex/plex_start.sh

[program:salt-minion]
command=/usr/bin/salt-minion -d